# cronjob/crontab
* * * * * echo hello
*/2  *  *   *   * /usr/local/bin/python3.8 /usr/src/app-staff/manage.py submit_newsletter 1>/dev/null 2>&1
# Mandatory blank line